# checkDrives.timer
[Unit]
Description=Run checkDrives once every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
