# checkDrives.service
[Unit]
Description=Run checkDrives

[Service]
Type=oneshot
ExecStart=/usr/local/bin/checkDrives.sh
